# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

#ipset start with system
ipset create chnroute hash:net
ipset restore -f /root/script/chnroute-ipset

exit 0
