# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

#ipset start with system
#ipset create chnroute hash:net
#ipset restore -f /etc/openvpn/chnroute-ipset


#关闭bird开机启动
/etc/init.d/bird stop
#/etc/init.d/bird disable

exit 0
